# scripts/rpa-auto-sync.service
[Unit]
Description=RPA Engine Auto-Sync
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/opt/rpa-engine/scripts/auto-sync.sh
Restart=always
RestartSec=10
User=root
WorkingDirectory=/opt/rpa-engine
Environment=HOME=/root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpa-auto-sync

[Install]
WantedBy=multi-user.target
